# pms/docker-mod/root/etc/s6-overlay/s6-rc.d/5-run-nginx/run
#!/usr/bin/with-contenv bash

echo "**** starting nginx ****"
/usr/sbin/nginx -g 'daemon off;'
