# pms/docker-mod/root/etc/s6-overlay/s6-rc.d/5-run-nginx/up
/etc/s6-overlay/s6-rc.d/4-setup-nginx/run