# worker/docker-mod/root/etc/s6-overlay/s6-rc.d/3-install-shim/up
/etc/s6-overlay/s6-rc.d/3-install-shim/run